# repo-upd.service
[Unit]
Description=Startup script repo-upd

[Service]
ExecStart=/usr/bin/repo-upd -m
